# steamlink.service.d/environment.conf
[Service]
Environment="PATH=/usr/osmc/bin:/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
